# dracon-system-guard.service
[Unit]
Description=Dracon System Guard - Proactive disk space monitoring and cleanup
Documentation=https://github.com/DraconDev/dracon-utilities
After=network.target

[Service]
Type=simple
StandardOutput=journal
StandardError=journal
# NixOS-optimized PATH. On non-NixOS distributions, customize or use
# PassEnvironment=PATH after running `systemctl --user import-environment PATH`.
Environment=PATH=%h/.local/bin:/run/wrappers/bin:%h/.local/share/flatpak/exports/bin:/var/lib/flatpak/exports/bin:%h/.nix-profile/bin:/nix/profile/bin:%h/.local/state/nix/profile/bin:/etc/profiles/per-user/%u/bin:/nix/var/nix/profiles/default/bin:/run/current-system/sw/bin
ExecStart=%h/.local/bin/dracon-system guard daemon
Restart=always
RestartSec=10
RestartPreventExitStatus=2 78

# Run as the user who owns the service
# For system-wide deployment, create a dedicated user
# User=dracon-guard
# Group=dracon-guard

# Resource limits
MemoryMax=250M
CPUQuota=20%
TasksMax=64

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.dracon %h/Dev %h/.local/state/dracon %h/.local/share/Trash %h/.cargo %h/.cache %h/.npm
PrivateTmp=true

# Kernel-level sandboxing (all certain — these daemons never use these privileges)
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
CapabilityBoundingSet=
RestrictNamespaces=true

# Syscall allowlist — returns EPERM instead of SIGSYS for safe discovery
SystemCallFilter=@system-service

[Install]
WantedBy=default.target
